# test.cron
00 16 * * * sh ~/home/xlqin/workdir/program/tonytam/spiderClub/CSSL_SOC/weibo.sh
