# test.cron
53 16 * * * cd /home/xlqin/workdir/program/tonytam/spiderClub/CSSL_SOC;sh weibo.sh
