# test.cron
59 23 * * * cd /home/xlqin/workdir/program/tonytam/spiderClub/CSSL_SOC;sh weibo.sh
