# test.cron
35 16 * * * sh ~/workdir/program/tonytam/spiderClub/CSSL_SOC/weibo.sh
